# server/systemd/udtun.service
[Unit]
Description=UDTUN UDP Tunneling Server
After=network.target
Wants=network.target
Before=network-online.target

[Service]
Type=notify
NotifyAccess=all
User=root
WorkingDirectory=/opt/udtun/server
ExecStartPre=/opt/udtun/scripts/setup_iptables.sh install
ExecStart=/usr/bin/python3 /opt/udtun/server/main.py
ExecStopPost=/opt/udtun/scripts/setup_iptables.sh remove
ExecReload=/bin/kill -HUP $MAINPID

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=0

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/udtun
ReadWritePaths=/run/udtun
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# Performance
LimitNOFILE=1048576
LimitCORE=infinity
LimitMEMLOCK=infinity

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=udtun

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
